# visualizations/visualizations/ds_scheduled-task-deleted.yml
title: DS_Scheduled task deleted
description: This rule detects scheduled task deletion.
status: Active
author: DS
logsource:
  product: Splunk
  service: splunkd
  category: Windows
detection:
  condition:
    sigma: "index=winlog EventCode=4699 AND user!=\"*$\" \n\n  -  table Task_Name\
      \ action user dest signature"
  level: critical
tags:
- T1053.005
- 'Scheduled Task/Job: Scheduled Task'
- Execution, Persistence, Privilege Escalation
references: .nan
notes: .nan
