# visualizations/visualizations/escu---linux-service-restarted---rule.yml
title: ESCU - Linux Service Restarted - Rule
description: This analytic looks for restarted or re-enable services in linux platform.
  This technique can be executed or performed using systemctl or service tool application.
  Adversaries may create or modify Windows services to repeatedly execute malicious
  payloads as part of persistence. When Windows boots up, it starts programs or applications
  called services that perform background system functions. Administrator may also
  create a legitimated service for a specific tool or normal application as part of
  task or automation, in this scenario it is suggested to look for the service path
  of the actual script or executable that register as service and who created the
  service for further verification.
status: Missing Data Source
author: DS
logsource:
  product: Splunk
  service: splunkd
  category: Linux
detection:
  condition:
    sigma: "\n  -  tstats `security_content_summariesonly` count min(_time) as firstTime\
      \ max(_time) as lastTime from datamodel=Endpoint.Processes where (Processes.process_name\
      \ IN (\"systemctl\", \"service\") OR Processes.process IN (\"*systemctl *\"\
      , \"*service *\")) Processes.process IN (\"*restart*\", \"*reload*\", \"*reenable*\"\
      ) by Processes.dest Processes.user Processes.parent_process_name Processes.process_name\
      \ Processes.process Processes.process_id Processes.parent_process_id Processes.process_guid\
      \ \n  -  `drop_dm_object_name(Processes)` \n  -  `security_content_ctime(firstTime)`\
      \ \n  -  `security_content_ctime(lastTime)` \n  -  `linux_service_restarted_filter`"
  level: critical
tags:
- T1053.006
- 'Scheduled Task/Job: Systemd Timers'
- Execution, Persistence, Privilege Escalation
references: Increase Linux logging to include command line data
notes: .nan
