# visualizations/visualizations/escu---scheduled-tasks-used-in-badrabbit-ransomware---rule.yml
title: ESCU - Scheduled tasks used in BadRabbit ransomware - Rule
description: 'WARNING, this detection has been marked deprecated by the Splunk Threat
  Research team, this means that it will no longer be maintained or supported. If
  you have any questions feel free to email us at: research@splunk.com. This search
  looks for flags passed to schtasks.exe on the command-line that indicate that task
  names related to the execution of Bad Rabbit ransomware were created or deleted.
  Deprecated because we already have a similar detection'
status: Inactive
author: DS
logsource:
  product: Splunk
  service: splunkd
  category: Windows
detection:
  condition:
    sigma: "\n  -  tstats `security_content_summariesonly` count min(_time) as firstTime\
      \ max(_time) as lastTime values(Processes.process) as process  from datamodel=Endpoint.Processes\
      \ where Processes.process_name=schtasks.exe (Processes.process= \"*create*\"\
      \  OR Processes.process= \"*delete*\") by Processes.parent_process Processes.process_name\
      \ Processes.user \n  -  `drop_dm_object_name(\"Processes\")` \n  -  `security_content_ctime(firstTime)`\n\
      \  - `security_content_ctime(lastTime)` \n  -  (process=*rhaegal* OR process=*drogon*\
      \ OR *viserion_*) \n  -  `scheduled_tasks_used_in_badrabbit_ransomware_filter`"
  level: critical
tags:
- T1053.005
- 'Scheduled Task/Job: Scheduled Task'
- Execution, Persistence, Privilege Escalation
references: .nan
notes: Deprecated by splunk
